# debian/percona-tokumxse-server.files
usr/sbin/mongod
usr/bin/thp-disable.sh
etc/init.d/mongod
etc/default/mongod
